# bt_rtk.service
[Unit]
Description=Bluetooth hciattach

[Service]
Type=simple
ExecStart=/usr/bin/rtk_hciattach -n -s 115200 ttyS0 rtk_h5
